# supervisord.conf
[supervisord]
user=root
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info
environment=HOME="/home/orb",USER="orb",ROS_HOME="/home/orb/.ros",XDG_RUNTIME_DIR="/tmp/runtime-orb",LD_LIBRARY_PATH="/home/orb/ORB_SLAM3/lib:/home/orb/colcon_ws/install/lib:/opt/ros/humble/lib",PATH="/opt/ros/humble/bin:/usr/bin"

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0777

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:slam_core]
command=/bin/bash -lc "
    export LD_LIBRARY_PATH='/home/orb/ORB_SLAM3/lib:/home/orb/colcon_ws/install/lib:/opt/ros/humble/lib:$LD_LIBRARY_PATH' &&
    export PATH='/opt/ros/humble/bin:$PATH' &&
    ros2 run orb_slam3_ros2_wrapper mono /home/orb/ORB_SLAM3/Vocabulary/ORBvoc.bin /home/orb/Database/real.yaml --ros-args -r __ns:=/"
user=orb
autostart=false
autorestart=false
startretries=3
stderr_logfile=/var/log/supervisor/slam_core.err.log
stdout_logfile=/var/log/supervisor/slam_core.out.log
environment=ROS_DOMAIN_ID="0",QT_X11_NO_MITSHM="1",LD_LIBRARY_PATH="/home/orb/ORB_SLAM3/lib:/home/orb/colcon_ws/install/lib:/opt/ros/humble/lib",PATH="/opt/ros/humble/bin:/usr/bin"
stopsignal=TERM
stopwaitsecs=15
stopasgroup=true
killasgroup=true

[program:relay]
command=/home/orb/TerraSLAM_relay/run_relay.sh
directory=/home/orb/TerraSLAM_relay
user=orb
autostart=false
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/relay.err.log
stdout_logfile=/var/log/supervisor/relay.out.log
environment=ROS_DOMAIN_ID="0"

# Folder mode (pre-recorded images)
[program:image_publisher_folder]
command=/bin/bash -l -c 'source /opt/ros/humble/setup.bash && source /home/orb/colcon_ws/install/setup.bash && cd /home/orb/Database && exec python3 image_publish.py "${VIDEO_FOLDER}"'
directory=/home/orb/Database
user=orb
autostart=false
autorestart=false
startretries=3
stderr_logfile=/var/log/supervisor/publisher_folder.err.log
stdout_logfile=/var/log/supervisor/publisher_folder.out.log
environment=VIDEO_FOLDER="Mill-video/",ROS_DOMAIN_ID="0"
stopsignal=TERM          ; сначала мягкий SIGTERM
stopwaitsecs=3          ; ждём 30 секунд перед SIGKILL
stopasgroup=true         ; убивать всю группу процессов
killasgroup=true         ; в том числе при SIGKILL

# RealSense mode (live camera) - disabled by default
[program:image_publisher_realsense]
command=/bin/bash -l -c "source /opt/ros/humble/setup.bash && source /home/orb/colcon_ws/install/setup.bash && exec python3 /home/orb/Database/realsense.py"
directory=/home/orb/Database
user=orb
autostart=false          ; ← Start manually when needed
autorestart=false
startretries=3
stderr_logfile=/var/log/supervisor/publisher_realsense.err.log
stdout_logfile=/var/log/supervisor/publisher_realsense.out.log
environment=ROS_DOMAIN_ID="0"

[program:rosbridge]
command=/bin/bash -l -c "source /opt/ros/humble/setup.bash && ros2 launch rosbridge_server rosbridge_websocket_launch.xml port:=9091 cors_origin:=* max_message_size:=10000000 fragment_timeout:=4000"
user=orb
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/rosbridge.out.log
stderr_logfile=/var/log/supervisor/rosbridge.err.log
environment=ROS_DOMAIN_ID="0"

[program:gps_bridge]
command=/bin/bash -l -c "source /opt/ros/humble/setup.bash && source /home/orb/colcon_ws/install/setup.bash && python3 /home/orb/TerraSLAM_relay/Serial/gps_bridge_node.py --ros-args -p protocol:=msp -p hw_type:=uart -p port:=/dev/ttyCH341USB0 -p baudrate:=115200"
directory=/home/orb/TerraSLAM_relay/Serial
user=orb
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/gps_bridge.err.log
stdout_logfile=/var/log/supervisor/gps_bridge.out.log

# 1. ROS2 System Manager — мозг, управляет supervisor, and with gateway
[program:system_manager]
command=/bin/bash -l -c "source /opt/ros/humble/setup.bash && source /home/orb/colcon_ws/install/setup.bash && cd /home/orb/TerraSLAM_relay && exec python3 system_manager.py"
directory=/home/orb/TerraSLAM_relay
user=orb
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/ros_system_manager.err.log
stdout_logfile=/var/log/supervisor/ros_system_manager.out.log
environment=ROS_DOMAIN_ID="0",PYTHONUNBUFFERED="1"

